# camera/camera.service
[Unit]
Description=Camera Service
After=network.target

[Service]
User=<your_username>
ExecStart=/usr/bin/go run /path/to/your/main.go
Restart=always

[Install]
WantedBy=multi-user.target